# ark.conf
[program:ark]
command = gunicorn ark.wsgi_test -b 127.0.0.1:8023
directory = /root/EgretProjects/ark
autostart = true
autorestart = true
